# automation/com.dealradarus.daily-automation.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.dealradarus.daily-automation</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/admin/projects/deal-aggregator-facebook/automation/startup-automation.sh</string>
    </array>
    
    <key>WorkingDirectory</key>
    <string>/Users/admin/projects/deal-aggregator-facebook</string>
    
    <key>RunAtLoad</key>
    <true/>
    
    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
        <key>NetworkState</key>
        <true/>
    </dict>
    
    <key>StandardOutPath</key>
    <string>/Users/admin/projects/deal-aggregator-facebook/automation/logs/launchd.log</string>
    
    <key>StandardErrorPath</key>
    <string>/Users/admin/projects/deal-aggregator-facebook/automation/logs/launchd-error.log</string>
    
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin:/opt/homebrew/bin</string>
        <key>NODE_ENV</key>
        <string>production</string>
    </dict>
    
    <!-- Start 60 seconds after login to ensure system stability -->
    <key>StartInterval</key>
    <integer>60</integer>
    
    <!-- Only run if network is available -->
    <key>LaunchOnlyOnce</key>
    <true/>
    
    <!-- Restart if the process crashes -->
    <key>ThrottleInterval</key>
    <integer>300</integer>
</dict>
</plist>